# systemd/accel-textfile.timer
[Unit]
Description=Run Intel iGPU + NPU telemetry collector every 15s

[Timer]
OnBootSec=30s
OnUnitActiveSec=15s
AccuracySec=1s

[Install]
WantedBy=timers.target
