# systemd/observability-onboard.timer
[Unit]
Description=Run observability onboarding every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=10s

[Install]
WantedBy=timers.target
